# Ubuntu Deploy Configs/Supervisor_conf/conf.d/ngs_server_worker.conf
[program:ngs_server_worker]
command=/home/hui/ngs_server/venv/bin/rq worker ngs-server-tasks
numprocs=3
process_name=%(program_name)s_%(process_num)02d
directory=/home/hui/ngs_server
user=hui
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

